# node/systemd/etcd-member.service.d/05-mount.conf
[Unit]
After=var-lib-etcd.mount
Requires=var-lib-etcd.mount

