# systemd/mnt-shared.mount
[Mount]
What=/dev/xvdg
Where=/mnt/shared
Type=ext4
