# node/systemd/var-lib-etcd.mount
[Mount]
What=/dev/xvdf
Where=/var/lib/etcd
Type=ext4
